# site5.conf
[program:site5] 
autorestart=true 
autostart=true 
numprocs=1 
numprocs_start=1 
startretries=0 
command=docker run -p 80:80 --rm --name site5 mysite2:php2 
startsecs=0